# rootfs/etc/s6-overlay/s6-rc.d/nginx-xray/run
#!/command/with-contenv sh
# shellcheck shell=sh

set -eou pipefail

# render-env reads {{VAR}} placeholders from the environment; allexport makes the
# values sourced from XRAY_RUNTIME_CONFIG (e.g. XRAY_XHTTP_PATH) visible to it.
set -o allexport

echo "Starting nginx-xray..."

# Load all xray variables from xray-config
# shellcheck disable=SC1090
. "${XRAY_RUNTIME_CONFIG}"

# s6 runs this service from the compiled servicedir, so PWD is not the source
# tree under /etc/s6-overlay. We render (and include) the configs from here, so
# nginx.conf's `include {{NGINX_DATA_DIR}}/...` must resolve to this same rendered
# copy rather than the untouched source file.
NGINX_DATA_DIR="${PWD}/data"
NGINX_CONFIG="${NGINX_DATA_DIR}/nginx.conf"

# Render the {{VAR}} placeholders in each nginx config in place. render-env
# discovers every placeholder automatically and aborts if any referenced
# variable is unset or empty; nginx's own $variables are left untouched.
# --in-place keeps each file's permissions and inode.
for nginx_template in "${NGINX_DATA_DIR}"/*.conf; do
    render-env --in-place "${nginx_template}"
done

exec s6-notifyoncheck -n 0 \
    nginx -g 'daemon off;' -c "${NGINX_CONFIG}"
