# rootfs/etc/s6-overlay/s6-rc.d/nginx-xray/run
#!/command/with-contenv sh
# shellcheck shell=sh

set -eou pipefail

# render-env reads {{VAR}} placeholders from the environment; allexport makes the
# values sourced from XRAY_RUNTIME_CONFIG (e.g. XRAY_XHTTP_PATH) visible to it.
set -o allexport

echo "Starting nginx-xray..."

# Load all xray variables from xray-config
# shellcheck disable=SC1090
. "${XRAY_RUNTIME_CONFIG}"

NGINX_CONFIG="${PWD}/data/nginx.conf"

# Render the {{VAR}} placeholders in each nginx config in place. render-env
# discovers every placeholder automatically and aborts if any referenced
# variable is unset or empty; nginx's own $variables are left untouched. Writing
# via a temp file + cat keeps the original file's permissions and inode.
for nginx_template in "${PWD}"/data/*.conf; do
    rendered=$(mktemp)
    render-env "${nginx_template}" > "${rendered}"
    cat "${rendered}" > "${nginx_template}"
    rm -f "${rendered}"
done

exec s6-notifyoncheck -n 0 \
    nginx -g 'daemon off;' -c "${NGINX_CONFIG}"
