# rootfs/etc/s6-overlay/s6-rc.d/xray-config/run
#!/command/with-contenv sh
# shellcheck shell=sh

set -eou pipefail

# Configure all variables to be exported to be readable by configgen
set -o allexport

echo "Starting xray config..."

mkdir -p "${GENERATED_CONFIG_PATH}"

# Check if our xray config file exists
if [ ! -f "${XRAY_CONFIG}" ]; then
    echo "${XRAY_CONFIG} not found, building a new config"

    XRAY_UUID=$(xray uuid)
    XRAY_REALITY_SHORT_ID=$(head -c 8 /dev/urandom | od -An -tx1 | tr -d ' \n')

    xray_keys=$(xray x25519)
    XRAY_REALITY_PRIVATE_KEY=$(echo "${xray_keys}" | grep 'PrivateKey' | cut -d ' ' -f2)
    XRAY_REALITY_PUBLIC_KEY=$(echo "${xray_keys}" | grep 'Password' | cut -d ' ' -f2)

    {
        echo "XRAY_UUID=${XRAY_UUID}"
        echo "XRAY_REALITY_SHORT_ID=${XRAY_REALITY_SHORT_ID}"
        echo "XRAY_REALITY_PRIVATE_KEY=${XRAY_REALITY_PRIVATE_KEY}"
        echo "XRAY_REALITY_PUBLIC_KEY=${XRAY_REALITY_PUBLIC_KEY}"
    } > "${XRAY_CONFIG}"
fi

{
    echo "XRAY_XHTTP_PATH=${XRAY_PATH}/xhttp"
} > "${XRAY_RUNTIME_CONFIG}"

. "${XRAY_CONFIG}"
. "${XRAY_RUNTIME_CONFIG}"

echo "Building xray client configs..."

# Get public IP
if [ -z "${PUBLIC_IP:+x}" ]; then
    PUBLIC_IP=$(wget -qO- api.ipify.org)
fi

# Pin the leaf certificate that nginx serves for XHTTP_DOMAIN. Clients use this
# (pinnedPeerCertSha256 / "pcs") instead of the removed "allowInsecure": they
# accept the exact cert regardless of SNI, but the connection stays
# cryptographically verified. The hash is sha256 over the DER of the leaf cert,
# which is what Xray-core compares against.
# The certbot/openssl-cert services (a dependency of this one) guarantee the
# cert is in place before we reach this point.
XHTTP_CERT_FILE="${CERTBOT_LIVE_CERT_PATH}/${XHTTP_DOMAIN}/fullchain.pem"
if [ ! -f "${XHTTP_CERT_FILE}" ]; then
    echo "Expected XHTTP certificate not found at ${XHTTP_CERT_FILE}" >&2
    exit 1
fi
XRAY_XHTTP_PINNED_CERT_SHA256=$(
    openssl x509 -in "${XHTTP_CERT_FILE}" -outform DER \
        | openssl dgst -sha256 \
        | sed 's/^.*= *//'
)
if [ ${#XRAY_XHTTP_PINNED_CERT_SHA256} -ne 64 ]; then
    echo "Failed to compute pinned cert SHA-256 for ${XHTTP_DOMAIN}" >&2
    exit 1
fi
echo "Pinned XHTTP cert SHA-256: ${XRAY_XHTTP_PINNED_CERT_SHA256}"

CLIENT_CONFIG_PATH="${GENERATED_CONFIG_PATH}/clients"
# Remove any old generated configs
if [ -d "${CLIENT_CONFIG_PATH}" ]; then
    rm -rf "${CLIENT_CONFIG_PATH}"
fi
mkdir -p "${CLIENT_CONFIG_PATH}"
mkdir -p "${CLIENT_CONFIG_PATH}/json"
mkdir -p "${CLIENT_CONFIG_PATH}/yaml"


# loop over all files in data/clients
# Render each template into a throwaway copy instead of editing it in place, so
# the shipped templates keep their {{placeholders}}. This keeps generation
# idempotent: the certbot deploy-hook can re-run this script after a cert
# renewal to refresh the pinned cert hash.
for client in /etc/xray/clients/*.yaml; do
    client_name=$(basename "${client}" .yaml)
    rendered=$(mktemp)
    sed \
        -e "s|{{PUBLIC_IP}}|${PUBLIC_IP}|g" \
        -e "s|{{XHTTP_DOMAIN}}|${XHTTP_DOMAIN}|g" \
        -e "s|{{REALITY_DOMAIN}}|${REALITY_DOMAIN}|g" \
        -e "s|{{XRAY_UUID}}|${XRAY_UUID}|g" \
        -e "s|{{XRAY_REALITY_PUBLIC_KEY}}|${XRAY_REALITY_PUBLIC_KEY}|g" \
        -e "s|{{XRAY_REALITY_SHORT_ID}}|${XRAY_REALITY_SHORT_ID}|g" \
        -e "s|{{XRAY_XHTTP_PATH}}|${XRAY_XHTTP_PATH}|g" \
        -e "s|{{XRAY_XHTTP_PINNED_CERT_SHA256}}|${XRAY_XHTTP_PINNED_CERT_SHA256}|g" \
        "${client}" \
        > "${rendered}"

    yq eval-all \
        '.inbounds += (select(fi == 1).inbounds) | select(fi == 0)' \
        "${rendered}" \
        /config/inbound.yaml \
        -oj \
        > "${CLIENT_CONFIG_PATH}/json/${client_name}.json"
    yq eval-all \
        '.inbounds += (select(fi == 1).inbounds) | select(fi == 0)' \
        "${rendered}" \
        /config/inbound.yaml \
        -oy \
        > "${CLIENT_CONFIG_PATH}/yaml/${client_name}.yaml"

    rm -f "${rendered}"
done

configgen /etc/xray/client_uri.yaml -o /etc/xray/client_uri.txt
