# rootfs/etc/s6-overlay/s6-rc.d/xray-config/run
#!/command/with-contenv sh
# shellcheck shell=sh

set -eou pipefail

# render-env reads {{VAR}} placeholders from the process environment, so every
# value we substitute (including the runtime-computed ones below: PUBLIC_IP, the
# VLESS strings, the pinned cert hashes, and everything sourced from the secret
# and runtime config files) must be exported. allexport does that for all
# assignments that follow.
set -o allexport

echo "Building xray configs..."

mkdir -p "${GENERATED_CONFIG_PATH}"

# Secrets are minted and persisted by the xray-secrets oneshot (a dependency of
# this service), so they are guaranteed to exist here. Source them first so the
# runtime config below can derive values from them.
# shellcheck disable=SC1090
. "${XRAY_SECRET_CONFIG}"


# 8 bytes (64 bits) of entropy as a 16-char hex string. Kept short on purpose:
# the resulting header name (X-Tk-<16hex> = 21 chars) must fit within nginx's
# fixed grpc_headers_hash bucket size of 64 bytes. A 32-hex token overflowed it
# and nginx refused to start (could not build grpc_headers_hash). 64 bits is
# ample for a secret-marker header name whose secrecy is the name itself.
XRAY_XFF_TOKEN="$(head -c 8 /dev/urandom | od -An -tx1 | tr -d ' \n')"
{
    echo "XRAY_XHTTP_PATH=${XRAY_PATH}/xhttp"
    echo "XRAY_TRUSTED_XFF_HEADER=X-Tk-${XRAY_XFF_TOKEN}"
} > "${XRAY_RUNTIME_CONFIG}"

# shellcheck disable=SC1090
. "${XRAY_RUNTIME_CONFIG}"

# VLESS Encryption ("Seed") strings, built from the persisted ML-KEM-768 keypair.
# Server "decryption" uses the seed with a "600s" resumption-ticket window;
# clients use the encapsulation key with "0rtt". The two must stay paired.
# They are consumed by render-env via the (allexport'd) environment rather than
# referenced directly here, hence the SC2034 suppression.
# shellcheck disable=SC2034
XRAY_VLESS_DECRYPTION="mlkem768x25519plus.native.600s.${XRAY_VLESS_ENC_SEED}"
# shellcheck disable=SC2034
XRAY_VLESS_ENCRYPTION="mlkem768x25519plus.native.0rtt.${XRAY_VLESS_ENC_CLIENT}"

# --- Server config -------------------------------------------------------------
# Render the {{VAR}} placeholders in server.yaml from the environment, then merge
# the result (read from stdin as the first document) with the user-provided
# outbounds.
XRAY_SERVER_CONFIG="/etc/xray/server.json"
render-env /etc/xray/server.yaml \
    | yq eval-all \
        '.outbounds += (select(fi == 1).outbounds) | select(fi == 0)' \
        - \
        /config/outbound.yaml \
        -o=json \
    > "${XRAY_SERVER_CONFIG}"

# --- Client configs ------------------------------------------------------------
echo "Building xray client configs..."

# Get public IP
if [ -z "${PUBLIC_IP:+x}" ]; then
    PUBLIC_IP=$(wget -qO- https://api.ipify.org)
    [ -n "${PUBLIC_IP}" ] || { echo "Failed to determine PUBLIC_IP" >&2; exit 1; }
fi

# Pin the leaf certificates nginx serves. Clients use these hashes
# (pinnedPeerCertSha256 / "pcs") instead of the removed "allowInsecure": they
# accept the exact cert regardless of trust store state, but the connection stays
# cryptographically verified. The hash is sha256 over the DER of the leaf cert,
# which is what Xray-core compares against.
# The certbot/openssl-cert services (dependencies of this one) guarantee the
# certs are in place before we reach this point.
compute_pinned_cert_sha256() {
    cert_file="$1"

    if [ ! -f "${cert_file}" ]; then
        echo "Expected certificate not found at ${cert_file}" >&2
        exit 1
    fi

    pinned_cert_sha256=$(
        openssl x509 -in "${cert_file}" -outform DER \
            | openssl dgst -sha256 \
            | sed 's/^.*= *//'
    )
    if [ ${#pinned_cert_sha256} -ne 64 ]; then
        echo "Failed to compute pinned cert SHA-256 for ${cert_file}" >&2
        exit 1
    fi

    echo "${pinned_cert_sha256}"
}

XRAY_XHTTP_PINNED_CERT_SHA256=$(compute_pinned_cert_sha256 "${CERTBOT_LIVE_CERT_PATH}/${XHTTP_DOMAIN}/fullchain.pem")
# nginx's catch-all default_server (the xhttp_self profile) reuses the
# REALITY_DOMAIN leaf, so pin that one.
XRAY_XHTTP_SELF_PINNED_CERT_SHA256=$(compute_pinned_cert_sha256 "${CERTBOT_LIVE_CERT_PATH}/${REALITY_DOMAIN}/fullchain.pem")
echo "Pinned XHTTP cert SHA-256: ${XRAY_XHTTP_PINNED_CERT_SHA256}"
echo "Pinned XHTTP self cert SHA-256: ${XRAY_XHTTP_SELF_PINNED_CERT_SHA256}"

CLIENT_CONFIG_PATH="${GENERATED_CONFIG_PATH}/clients"
# Remove any old generated configs
if [ -d "${CLIENT_CONFIG_PATH}" ]; then
    rm -rf "${CLIENT_CONFIG_PATH}"
fi
mkdir -p "${CLIENT_CONFIG_PATH}"
mkdir -p "${CLIENT_CONFIG_PATH}/json"
mkdir -p "${CLIENT_CONFIG_PATH}/yaml"


# loop over all files in data/clients
# render-env reads each template and writes to a throwaway copy without touching
# the source file, so the shipped templates keep their {{placeholders}}. This
# keeps generation idempotent: the certbot deploy-hook can re-run this script
# after a cert renewal to refresh the pinned cert hash. The rendered copy is
# reused for both the JSON and YAML outputs.
for client in /etc/xray/clients/*.yaml; do
    client_name=$(basename "${client}" .yaml)
    rendered=$(mktemp)
    render-env "${client}" > "${rendered}"
    yq eval-all \
        '.inbounds += (select(fi == 1).inbounds) | select(fi == 0)' \
        "${rendered}" \
        /config/inbound.yaml \
        -o=json \
        > "${CLIENT_CONFIG_PATH}/json/${client_name}.json"
    yq eval-all \
        '.inbounds += (select(fi == 1).inbounds) | select(fi == 0)' \
        "${rendered}" \
        /config/inbound.yaml \
        -o=yaml \
        > "${CLIENT_CONFIG_PATH}/yaml/${client_name}.yaml"
    rm -f "${rendered}"
done

# --- Client share links + QR codes --------------------------------------------
# Replaces the former `configgen` binary. Render the {{placeholders}} with
# render-env, convert to compact JSON with yq, build + percent-encode the share
# links with render-client-uri, then render a QR code per link.
echo "Building xray client share links..."

render-env /etc/xray/client_uri.yaml \
    | yq -o=json -I=0 \
    | render-client-uri \
    > /etc/xray/client_uri.txt

while IFS= read -r uri; do
    # print the uri in the format of Name:\n uri\n
    printf '%s:\n%s\n' "${uri##*#}" "${uri}"
    # Some profiles might be too large to encode, so we ignore the error
    qrencode -t UTF8i -m 1 "${uri}" || true
    echo
done < /etc/xray/client_uri.txt
