# rootfs/etc/s6-overlay/s6-rc.d/xray-config/run
#!/command/with-contenv sh
# shellcheck shell=sh

set -eou pipefail

# Export everything: configgen reads {{VARS}} from the environment, and the
# values are reused across the sed/yq rendering below.
set -o allexport

echo "Building xray configs..."

mkdir -p "${GENERATED_CONFIG_PATH}"

# Runtime (non-secret) config derived from env.
{
    echo "XRAY_XHTTP_PATH=${XRAY_PATH}/xhttp"
} > "${XRAY_RUNTIME_CONFIG}"

# Secrets are minted and persisted by the xray-secrets oneshot (a dependency of
# this service), so they are guaranteed to exist here.
# shellcheck disable=SC1090
. "${XRAY_SECRET_CONFIG}"
# shellcheck disable=SC1090
. "${XRAY_RUNTIME_CONFIG}"

# VLESS Encryption ("Seed") strings, built from the persisted ML-KEM-768 keypair.
# Server "decryption" uses the seed with a "600s" resumption-ticket window;
# clients use the encapsulation key with "0rtt". The two must stay paired.
XRAY_VLESS_DECRYPTION="mlkem768x25519plus.native.600s.${XRAY_VLESS_ENC_SEED}"
XRAY_VLESS_ENCRYPTION="mlkem768x25519plus.native.0rtt.${XRAY_VLESS_ENC_CLIENT}"

# --- Server config -------------------------------------------------------------
# Merge server.yaml with the user-provided outbounds, then render placeholders.
XRAY_SERVER_CONFIG="/etc/xray/server.json"
yq eval-all \
    '.outbounds += (select(fi == 1).outbounds) | select(fi == 0)' \
    /etc/xray/server.yaml \
    /config/outbound.yaml \
    -oj \
    > "${XRAY_SERVER_CONFIG}"
sed -i \
    -e "s|{{CERTBOT_LIVE_CERT_PATH}}|${CERTBOT_LIVE_CERT_PATH}|g" \
    -e "s|{{REALITY_DOMAIN}}|${REALITY_DOMAIN}|g" \
    -e "s|{{XHTTP_DOMAIN}}|${XHTTP_DOMAIN}|g" \
    -e "s|{{XRAY_XHTTP_PATH}}|${XRAY_XHTTP_PATH}|g" \
    -e "s|{{XRAY_UUID}}|${XRAY_UUID}|g" \
    -e "s|{{XRAY_REALITY_PRIVATE_KEY}}|${XRAY_REALITY_PRIVATE_KEY}|g" \
    -e "s|{{XRAY_REALITY_SHORT_ID}}|${XRAY_REALITY_SHORT_ID}|g" \
    -e "s|{{XRAY_REALITY_MLDSA65_SEED}}|${XRAY_REALITY_MLDSA65_SEED}|g" \
    -e "s|{{XRAY_VLESS_DECRYPTION}}|${XRAY_VLESS_DECRYPTION}|g" \
    "${XRAY_SERVER_CONFIG}"

# --- Client configs ------------------------------------------------------------
echo "Building xray client configs..."

if [ -z "${XHTTP_SELF_DOMAIN:-}" ]; then
    echo "XHTTP_SELF_DOMAIN must be set for the xhttp_self client profile" >&2
    exit 1
fi

# Get public IP
if [ -z "${PUBLIC_IP:+x}" ]; then
    PUBLIC_IP=$(wget -qO- api.ipify.org)
fi

# Pin the leaf certificates nginx serves. Clients use these hashes
# (pinnedPeerCertSha256 / "pcs") instead of the removed "allowInsecure": they
# accept the exact cert regardless of trust store state, but the connection stays
# cryptographically verified. The hash is sha256 over the DER of the leaf cert,
# which is what Xray-core compares against.
# The certbot/openssl-cert services (dependencies of this one) guarantee the
# certs are in place before we reach this point.
compute_pinned_cert_sha256() {
    domain="$1"
    cert_file="${CERTBOT_LIVE_CERT_PATH}/${domain}/fullchain.pem"

    if [ ! -f "${cert_file}" ]; then
        echo "Expected certificate not found at ${cert_file}" >&2
        exit 1
    fi

    pinned_cert_sha256=$(
        openssl x509 -in "${cert_file}" -outform DER \
            | openssl dgst -sha256 \
            | sed 's/^.*= *//'
    )
    if [ ${#pinned_cert_sha256} -ne 64 ]; then
        echo "Failed to compute pinned cert SHA-256 for ${domain}" >&2
        exit 1
    fi

    echo "${pinned_cert_sha256}"
}

XRAY_XHTTP_PINNED_CERT_SHA256=$(compute_pinned_cert_sha256 "${XHTTP_DOMAIN}")
XRAY_XHTTP_SELF_PINNED_CERT_SHA256=$(compute_pinned_cert_sha256 "${XHTTP_SELF_DOMAIN}")
echo "Pinned XHTTP cert SHA-256: ${XRAY_XHTTP_PINNED_CERT_SHA256}"
echo "Pinned XHTTP self cert SHA-256: ${XRAY_XHTTP_SELF_PINNED_CERT_SHA256}"

CLIENT_CONFIG_PATH="${GENERATED_CONFIG_PATH}/clients"
# Remove any old generated configs
if [ -d "${CLIENT_CONFIG_PATH}" ]; then
    rm -rf "${CLIENT_CONFIG_PATH}"
fi
mkdir -p "${CLIENT_CONFIG_PATH}"
mkdir -p "${CLIENT_CONFIG_PATH}/json"
mkdir -p "${CLIENT_CONFIG_PATH}/yaml"


# loop over all files in data/clients
# Render each template into a throwaway copy instead of editing it in place, so
# the shipped templates keep their {{placeholders}}. This keeps generation
# idempotent: the certbot deploy-hook can re-run this script after a cert
# renewal to refresh the pinned cert hash.
for client in /etc/xray/clients/*.yaml; do
    client_name=$(basename "${client}" .yaml)
    rendered=$(mktemp)
    sed \
        -e "s|{{PUBLIC_IP}}|${PUBLIC_IP}|g" \
        -e "s|{{XHTTP_DOMAIN}}|${XHTTP_DOMAIN}|g" \
        -e "s|{{XHTTP_SELF_DOMAIN}}|${XHTTP_SELF_DOMAIN}|g" \
        -e "s|{{REALITY_DOMAIN}}|${REALITY_DOMAIN}|g" \
        -e "s|{{XRAY_UUID}}|${XRAY_UUID}|g" \
        -e "s|{{XRAY_REALITY_PUBLIC_KEY}}|${XRAY_REALITY_PUBLIC_KEY}|g" \
        -e "s|{{XRAY_REALITY_MLDSA65_VERIFY}}|${XRAY_REALITY_MLDSA65_VERIFY}|g" \
        -e "s|{{XRAY_REALITY_SHORT_ID}}|${XRAY_REALITY_SHORT_ID}|g" \
        -e "s|{{XRAY_XHTTP_PATH}}|${XRAY_XHTTP_PATH}|g" \
        -e "s|{{XRAY_XHTTP_PINNED_CERT_SHA256}}|${XRAY_XHTTP_PINNED_CERT_SHA256}|g" \
        -e "s|{{XRAY_XHTTP_SELF_PINNED_CERT_SHA256}}|${XRAY_XHTTP_SELF_PINNED_CERT_SHA256}|g" \
        -e "s|{{XRAY_VLESS_ENCRYPTION}}|${XRAY_VLESS_ENCRYPTION}|g" \
        "${client}" \
        > "${rendered}"

    yq eval-all \
        '.inbounds += (select(fi == 1).inbounds) | select(fi == 0)' \
        "${rendered}" \
        /config/inbound.yaml \
        -oj \
        > "${CLIENT_CONFIG_PATH}/json/${client_name}.json"
    yq eval-all \
        '.inbounds += (select(fi == 1).inbounds) | select(fi == 0)' \
        "${rendered}" \
        /config/inbound.yaml \
        -oy \
        > "${CLIENT_CONFIG_PATH}/yaml/${client_name}.yaml"

    rm -f "${rendered}"
done

configgen /etc/xray/client_uri.yaml -o /etc/xray/client_uri.txt
