# rootfs/etc/s6-overlay/s6-rc.d/xray/run
#!/command/with-contenv sh
# shellcheck shell=sh

set -eou pipefail

echo "Starting xray..."

# /etc/xray/server.json is rendered by the xray-config oneshot (reached
# transitively via nginx-xray), so this service only runs the binary.
exec xray run -c /etc/xray/server.json
